# deploy/proof-factory-hard.timer
[Unit]
Description=Proof Factory - event-gated R(5,5) research review

[Timer]
OnCalendar=*-*-* 00,04,08,12,16,20:00:00 UTC
Persistent=true
AccuracySec=1min
Unit=proof-factory-hard.service

[Install]
WantedBy=timers.target
