# deploy/proof-factory-easy.service
[Unit]
Description=Proof Factory - Sol principal with Terra delegate easier-problem pass
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/proof-factory
Environment=HOME=/root
Environment=PATH=/root/proof-factory/.lab-venv/bin:/root/.cache/proof-factory/lean/elan/bin:/root/.elan/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=ELAN_HOME=/root/.cache/proof-factory/lean/elan
Environment=PROOF_FACTORY_CACHE_DIR=/root/.cache/proof-factory
Environment=PROOF_FACTORY_ROOT=/root/proof-factory
Environment=CODEX_BIN=/root/.local/bin/codex
Environment=PROOF_CODEX_SANDBOX=danger-full-access
Environment=PROOF_EASY_DELEGATE_TIMEOUT_SEC=600
EnvironmentFile=-/root/project-factory/.env
Environment=PHOENIX_ENABLED=true
Environment=PHOENIX_ENDPOINT=https://app.phoenix.arize.com/s/ctkrug4501
Environment=PHOENIX_PROJECT=proof-factory
Environment=PHOENIX_CAPTURE_TEXT=true
ExecStart=/root/proof-factory/.venv/bin/python -m proof_factory tick --lane easy
Nice=15
CPUWeight=2
IOWeight=10
MemoryHigh=700M
MemoryMax=1200M
OOMScoreAdjust=700
TimeoutStartSec=2h

[Install]
WantedBy=multi-user.target
