# deploy/proof-factory-hard.timer
[Unit]
Description=Proof Factory - twice-hourly around-the-clock R(5,5) research

[Timer]
OnCalendar=*-*-* *:00,30:00 UTC
Persistent=true
AccuracySec=1min
Unit=proof-factory-hard.service

[Install]
WantedBy=timers.target
