# deploy/proof-factory-scout.service
[Unit]
Description=Proof Factory - broad academic contribution scout
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/proof-factory
EnvironmentFile=/etc/proof-factory/proof-factory.env
Environment=HOME=/root
ExecStartPre=/root/proof-factory/.venv/bin/python -m proof_factory.config validate --profile scout
ExecStart=/root/proof-factory/.venv/bin/python -m proof_factory scout
ExecStartPost=/bin/systemctl start --no-block proof-factory-publish.service
Nice=18
CPUWeight=1
IOWeight=5
MemoryHigh=500M
MemoryMax=900M
OOMScoreAdjust=750
TimeoutStartSec=75min
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/proof-factory/data /root/proof-factory/state /root/proof-factory/site /root/proof-factory/research /root/.codex
InaccessiblePaths=/root/project-factory /root/.ssh /root/.config/gh
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
