# deploy/proof-factory-capacity-guard.service
[Unit]
Description=Proof Factory - disposable residue cleanup and capacity check
After=local-fs.target

[Service]
Type=oneshot
WorkingDirectory=/root/proof-factory
EnvironmentFile=/etc/proof-factory/proof-factory.env
Environment=HOME=/root
ExecStartPre=/root/proof-factory/.venv/bin/python -m proof_factory.config validate --profile capacity
ExecStart=/root/proof-factory/.venv/bin/python -m proof_factory capacity-guard --lane easy
Nice=19
CPUWeight=1
IOWeight=1
MemoryMax=256M
