# deploy/proof-factory-hard.timer
[Unit]
Description=Proof Factory - event-gated R(5,5) research review

[Timer]
OnBootSec=5min
OnUnitInactiveSec=15min
Persistent=true
AccuracySec=1min
Unit=proof-factory-hard.service

[Install]
WantedBy=timers.target
