# deploy/proof-factory-runtime-sync.service
[Unit]
Description=Proof Factory - publish live public runtime snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/proof-factory
EnvironmentFile=/etc/proof-factory/proof-factory.env
Environment=HOME=/root
EnvironmentFile=-/root/project-factory/.env
ExecStartPre=/root/proof-factory/.venv/bin/python -m proof_factory.config validate --profile runtime
ExecStart=/root/proof-factory/.venv/bin/python /root/proof-factory/scripts/publish_runtime.py
Nice=15
CPUWeight=2
IOWeight=10
MemoryMax=128M
TimeoutStartSec=1min
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/proof-factory/state
NoNewPrivileges=true
