# deploy/proof-factory-lab.service
[Unit]
Description=Proof Factory - checkpointed continuation experiment segment worker
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/proof-factory
EnvironmentFile=/etc/proof-factory/proof-factory.env
Environment=HOME=/root
Environment=PATH=/root/proof-factory/.lab-venv/bin:/root/.elan/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStartPre=/root/proof-factory/.venv/bin/python -m proof_factory.config validate --profile lab
ExecStart=/root/proof-factory/.venv/bin/python -m proof_factory lab-worker --drain
Nice=19
CPUWeight=1
CPUQuota=100%
IOWeight=5
MemoryHigh=3400M
MemoryMax=3600M
OOMScoreAdjust=800
TasksMax=128
TimeoutStartSec=25h
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/proof-factory/research /root/proof-factory/state
InaccessiblePaths=/root/project-factory /root/.ssh /root/.codex /root/.config/gh
IPAddressDeny=any
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
